# newspaper.timer
[Unit]
Description=Pull latest paper every day at 7am

[Timer]
OnCalendar=*-*-* 7:00:00
Persistent=true

[Install]
WantedBy=timers.target